# agents/gollum/com.trainingrun.gollum.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trainingrun.gollum</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/davidsolomon/Desktop/TR2/agents/gollum/scout.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TR_REPO_PATH</key>
        <string>/Users/davidsolomon/Desktop/TR2</string>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/davidsolomon/Desktop/TR2/agents/gollum</string>

    <key>StandardOutPath</key>
    <string>/Users/davidsolomon/Desktop/TR2/agents/gollum/scout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davidsolomon/Desktop/TR2/agents/gollum/scout.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
